# moxa-upload/moxa-archive-upload.timer
[Unit]
Description=Nightly upload of weather DB to radxa

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
